# packaging/ironcache.service
# SPDX-License-Identifier: MIT OR Apache-2.0
#
# Hardened systemd unit for the IronCache server (#122).
# SCAFFOLD: templated and inert until the first release. See
# docs/design/PACKAGING.md and packaging/README.md.
#
# Install to /etc/systemd/system/ironcache.service, then:
#   systemctl daemon-reload && systemctl enable --now ironcache
# Verify the sandbox with:
#   systemd-analyze security ironcache.service
#
# This is the in-memory default. For a persisted-data deployment, uncomment
# StateDirectory + ReadWritePaths so the data dir is writable under the
# otherwise read-only root filesystem.

[Unit]
Description=IronCache server
Documentation=https://github.com/OWNER/REPO
After=network-online.target
Wants=network-online.target
# Socket activation (#389): when ironcache.socket is enabled, systemd creates the
# listening socket FIRST and passes its fd to this service, which ADOPTS it instead
# of self-binding -- so the listen queue survives an upgrade restart with no
# connection-refused window. HARMLESS when the socket unit is absent: systemd
# ignores a Wants/After of a non-installed unit, and the server self-binds exactly
# as before. Wants (not Requires) keeps this fail-open: if the socket unit failed,
# the service still starts and self-binds rather than being blocked.
Wants=ironcache.socket
After=ironcache.socket

[Service]
Type=simple
# `--metrics-addr` exposes /metrics, /livez, /readyz on loopback. This is REQUIRED for
# `ironcache upgrade`: its health gate probes /readyz + scrapes /metrics (ironcache_uptime_seconds)
# at 127.0.0.1:9121 by default, and pre-flights that endpoint before swapping. Without it every
# upgrade's health gate would be connection-refused and would auto-roll-back a healthy binary.
ExecStart=/usr/local/bin/ironcache server --metrics-addr 127.0.0.1:9121
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2s

# --- Identity: transient, allocated-per-start UID/GID, no persistent account.
DynamicUser=yes

# --- No privilege gain through execve (setuid/setgid + file caps neutralized).
NoNewPrivileges=yes

# --- Drop every capability; the daemon needs none (binds an unprivileged port).
CapabilityBoundingSet=
AmbientCapabilities=

# --- Filesystem: entire tree read-only except explicitly allowed paths.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

# --- Memory: forbid writable-and-executable mappings (W^X).
MemoryDenyWriteExecute=yes

# --- Sockets: allow only the families the server actually uses.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# --- Additional defense-in-depth (consistent with the eight load-bearing
#     directives above; tightened further as the binary's syscall set settles).
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
UMask=0077

# --- Persisted-data deployments only: uncomment to grant a writable state dir
#     under the strict read-only root.
# StateDirectory=ironcache
# ReadWritePaths=/var/lib/ironcache

[Install]
WantedBy=multi-user.target
